# systemd/venus-export-limiter.service
[Unit]
Description=Venus OS Export Limiter v2 (Minimal, no logging)
After=dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /data/venus-export-limiter-v2/limit-control-v2.py
Restart=always

[Install]
WantedBy=multi-user.target
